# raspberry_code/cep-agent.service
[Unit]
Description=CEP Agent — Raspberry Pi (ponte MQTT, agendador diário/mensal, FastAPI)
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=mrmorais
WorkingDirectory=/home/mrmorais/MyBookRegister_by_JMorais/raspberry_code/cep_agent
EnvironmentFile=/home/mrmorais/MyBookRegister_by_JMorais/raspberry_code/cep_agent/.env
ExecStart=/home/mrmorais/MyBookRegister_by_JMorais/raspberry_code/venv_cep_agent/bin/uvicorn app:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
