# raspberry_code/lampada-stats.service
[Unit]
Description=Receptor de eventos da lâmpada (Flask) — Raspberry Pi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mrmorais
WorkingDirectory=/home/mrmorais/MyBookRegister_by_JMorais/raspberry_code
ExecStart=/home/mrmorais/MyBookRegister_by_JMorais/raspberry_code/venv_lampada/bin/gunicorn -w 1 -b 0.0.0.0:5000 lampada_stats:app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
